# config/config/systemd/user/low-battery-notifier.service
[Unit]
Description=Notify low battery alert

[Service]
Type=oneshot
ExecStart=/home/masawada/.dotfiles/bin/low-battery-notifier
